# packages/net-fs/samba/files/systemd/nmbd.service
[Unit]
Description=Samba NetBIOS over IP name server

[Service]
ExecStart=/usr/sbin/nmbd --foreground --log-stdout
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
