# packages/net-fs/samba/files/systemd/smbd.service
[Unit]
Description=Samba SMB/CIFS server

[Service]
ExecStart=/usr/sbin/smbd --foreground --log-stdout
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
