# packages/net-fs/samba/files/systemd/swat@.service
[Unit]
Description=Swat (Samba web configuration interface) daemon

[Service]
ExecStart=/usr/sbin/swat
StandardInput=socket
