# packages/net-fs/samba/files/systemd/winbindd.service
[Unit]
Description=Samba Name Service Switch server

[Service]
ExecStart=/usr/sbin/winbindd --foreground --log-stdout
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
